# rootfs/usr/lib/systemd/system/local-apt-repository.path
[Path]
PathChanged=/srv/local-apt-repository

[Install]
WantedBy=paths.target
