# rootfs/usr/lib/systemd/system/local-apt-repository.service
[Unit]
Description=Local apt repository recreation

[Service]
Type=oneshot
ExecStart=/usr/lib/local-apt-repository/rebuild
